# rootfs/etc/services.d/microsocks/run
#!/usr/bin/execlineb -P
s6-setuidgid microsocks
/usr/local/bin/microsocks
